# files/home/.config/systemd/user/borgmatic-backup.service
[Unit]
Description=Borgmatic backup

[Service]
Type=oneshot
ExecStart=/home/zayon/.local/bin/borgmatic --verbosity 1 --log-file-verbosity 1

