# files/home/.config/systemd/user/borgmatic-backup.timer
[Unit]
Description=Run borgmatic backup daily at 4 am

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
